# supervisor/backend.conf
[program:voice-text-processor-backend]
command=%(here)s/../scripts/start_backend.sh
    --bind 0.0.0.0:3082
    --workers 1
    --threads 2
    --timeout 1800
    --keep-alive 300
    --access-logfile -
    --error-logfile -
    --log-level info
    --reload
    --worker-class sync
    app:app
directory=%(here)s/../backend
environment=PORT=3082,BACKEND_URL=http://localhost:3082,AI_ENGINE=ollama,OLLAMA_TIMEOUT=1800,WHISPER_MODEL_FIXED=SoybeanMilk/faster-whisper-Breeze-ASR-25,OLLAMA_MODEL_FIXED=gpt-oss:20b,MAX_UPLOAD_SIZE=15360,WHISPER_BACKEND=faster_whisper,WHISPER_DEVICE=cuda,WHISPER_COMPUTE_TYPE=float16,CORS_ORIGINS=*
stdout_logfile=%(here)s/../logs/backend.log
stderr_logfile=%(here)s/../logs/backend_error.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
autostart=false
autorestart=true
startsecs=5
stopwaitsecs=10
stopasgroup=true
killasgroup=true
redirect_stderr=false
