# supervisor/vllm.conf
[program:voice-text-processor-vllm]
command=%(here)s/../scripts/start_vllm.sh
directory=%(here)s/..
environment=VLLM_MODEL=google/gemma-3-27b-it,VLLM_PORT=8000,VLLM_MAX_MODEL_LEN=65536,VLLM_GPU_MEM=0.85,HF_TOKEN=%(ENV_HF_TOKEN)s
stdout_logfile=%(here)s/../logs/vllm.log
stderr_logfile=%(here)s/../logs/vllm_error.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
autostart=false
autorestart=true
startsecs=60
stopwaitsecs=60
redirect_stderr=false
